# source/bbb-rootfs/etc/inittab
ttyO0::respawn:/sbin/getty -L ttyO0 115200 vt100
